# release/src-rt-6.x.4708/router/nfs-utils/systemd/nfs-blkmap.service
[Unit]
Description=pNFS block layout mapping daemon
DefaultDependencies=no
Conflicts=umount.target
After=var-lib-nfs-rpc_pipefs.mount
Requires=var-lib-nfs-rpc_pipefs.mount

PartOf=nfs-utils.service

[Service]
Type=forking
PIDFile=/var/run/blkmapd.pid
ExecStart=/usr/sbin/blkmapd

[Install]
WantedBy=nfs-client.target
